# gateway-bot.service
[Unit]
Description=Gateway Security Analyzer Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/gateway-security-analyzer-bot
Environment=PATH=/usr/bin:/usr/local/bin
EnvironmentFile=/path/to/gateway-security-analyzer-bot/.env
ExecStart=/usr/bin/python3 /path/to/gateway-security-analyzer-bot/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target